# Server/src/bin/openas2.service
[Unit]
Description=OpenAS2 Server
After=network.target

[Service]
Type=forking
User=root
Group=root

Environment=OPENAS2_AS_DAEMON=true
Environment=OPENAS2_BASE=/opt/OpenAS2
Environment=OPENAS2_PID=/opt/OpenAS2/bin/OpenAS2.pid
PIDFile=/opt/OpenAS2/bin/OpenAS2.pid
WorkingDirectory=/opt/OpenAS2
ExecStart=/opt/OpenAS2/bin/start-openas2.sh
SuccessExitStatus=0 143

[Install]
WantedBy=multi-user.target
